# data/supervisord/conf.d/gcs_cloud_integration_service_process.conf
[program:gcs-cloud-integration-service]
command=/bin/echo "Hello, Supervisord 2!"
user=gcs-cloud
autostart=true
autorestart=true
stdout_logfile=/var/log/gcs-cloud-integration-service/stdout.log
stderr_logfile=/var/log/gcs-cloud-integration-service/stderr.log