# data/supervisord/conf.d/gcs_broker_lora_rtk_process.conf
[program:gcs-broker-lora-rtk]
command=/bin/echo "Hello, Supervisord 1!"
user=gcs
autostart=true
autorestart=true
stdout_logfile=/var/log/gcs-broker-lora-rtk/stdout.log
stderr_logfile=/var/log/gcs-broker-lora-rtk/stderr.log