# data/supervisord/conf.d/gcs_ui_backend_process.conf
[program:gcs-ui-backend]
command=/bin/echo "Hello, Supervisord 5!"
user=gcs
autostart=true
autorestart=true
stdout_logfile=/var/log/gcs-ui-backend/stdout.log
stderr_logfile=/var/log/gcs-ui-backend/stderr.log