# data/supervisord/conf.d/gcs_broker_lora_rtk_process.conf
[program:process]
directory=/home/sergiy/program/
command=/usr/bin/php process.php
user=sergiy
autostart=true
autorestart=true