# data/supervisord/conf.d/gcs_connection_manager_process.conf
[program:process]
command=/bin/echo "Hello, Supervisord 3!"
user=gcs
autostart=true
autorestart=true
stdout_logfile=/var/log/gcs-connection-manager/stdout.log
stderr_logfile=/var/log/gcs-connection-manager/stderr.log