# data/supervisord/conf.d/opensearch_process.conf
[program:process]
command=/opt/opensearch/opensearch-tar-install.sh -Ecluster.name=opensearch-cluster -Enode.name=opensearch-node1 -Ehttp.host=0.0.0.0 -Ediscovery.type=single-node &
user=opensearch
autostart=true
autorestart=true
stdout_logfile=/var/log/opensearch/stdout.log
stderr_logfile=/var/log/opensearch/stderr.log