# data/supervisord/conf.d/opensearch_dashboards_process.conf
[program:opensearch-dashboards]
command=/opt/opensearch-dashboards/bin/opensearch-dashboards
user=opensearch
autostart=true
autorestart=true
stdout_logfile=/var/log/opensearch-dashboards/stdout.log
stderr_logfile=/var/log/opensearch-dashboards/stderr.log