# data/supervisord/conf.d/opensearch_process.conf
[program:process]
command=/opt/opensearch/bin/opensearch
user=opensearch
autostart=true
autorestart=true
stdout_logfile=/var/log/opensearch/stdout.log
stderr_logfile=/var/log/opensearch/stderr.log